# systemd/surveil-signal-outcome.service
[Unit]
Description=Surveil investment signal outcome backfill
After=network-online.target surveil-db-init.service
Wants=network-online.target

[Service]
Type=oneshot
User=surveil
WorkingDirectory=/opt/surveil
EnvironmentFile=-/opt/surveil/.env
ExecStart=/opt/surveil/.venv/bin/python /opt/surveil/scripts/signal_outcome_update.py --days 45
StandardOutput=append:/opt/surveil/logs/signal-outcome.log
StandardError=append:/opt/surveil/logs/signal-outcome.err.log
