# systemd/surveil-portfolio-import.service
[Unit]
Description=Surveil portfolio import
After=network-online.target surveil-db-init.service
Wants=network-online.target

[Service]
Type=oneshot
User=surveil
WorkingDirectory=/opt/surveil
EnvironmentFile=-/opt/surveil/.env
ExecStart=/opt/surveil/.venv/bin/python /opt/surveil/scripts/portfolio_import.py
StandardOutput=append:/opt/surveil/logs/portfolio-import.log
StandardError=append:/opt/surveil/logs/portfolio-import.err.log
