# systemd/surveil-sina-stock-news.timer
[Unit]
Description=Run Surveil Sina per-stock news batch every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true
Unit=surveil-sina-stock-news.service

[Install]
WantedBy=timers.target
