# systemd/surveil-ifind-smoke.service
[Unit]
Description=Surveil iFinD smoke test
After=network-online.target surveil-db-init.service
Wants=network-online.target

[Service]
Type=oneshot
User=surveil
WorkingDirectory=/opt/surveil
EnvironmentFile=-/opt/surveil/.env
ExecStart=/opt/surveil/.venv/bin/python /opt/surveil/scripts/smoke_ifind.py
StandardOutput=append:/opt/surveil/logs/ifind-smoke.log
StandardError=append:/opt/surveil/logs/ifind-smoke.err.log
