# launchd/io.github.surveil.x-stream.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.github.surveil.x-stream</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON_BIN__</string>
    <string>__SURVEIL_ROOT__/scripts/x_stream.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__SURVEIL_ROOT__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__SURVEIL_ROOT__/logs/x_stream.out.log</string>

  <key>StandardErrorPath</key>
  <string>__SURVEIL_ROOT__/logs/x_stream.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
